# templates/default/carbon.upstart.erb
#!upstart
# http://upstart.ubuntu.com/cookbook/
#
description carbon-<%= @name %>

start on runlevel [2345]
stop on shutdown

respawn
respawn limit 99 5

setuid <%= @user %>

pre-start script
  <% @instances.each do |instance| -%>
  LOGDIR="<%= @storage_dir %>/log/carbon-<%= @name %>/carbon-<%= @name %>-<%= instance.first %>"
  [ -d $LOGDIR ] || mkdir -p $LOGDIR
  chown -R <%= @user %> $LOGDIR
  <% end -%>
end script

script
  <% @instances.each do |instance| -%>
  <%= @base_dir %>/bin/carbon-<%= @name %>.py --debug --instance=<%= instance.first %> start >> <%= @storage_dir %>/log/carbon-<%= @name %>/carbon-<%= @name %>-<%= instance.first %>/console.log 2>&1
  <% end -%>
end script
