# templates/default/carbon-cache.conf.erb
description "Graphite Carbon Cache Daemon"

start on runlevel [2345]
stop on runlevel [!2345]

expect daemon
respawn

env PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

exec sudo -u www-data /opt/graphite/bin/carbon-cache.py start
