# templates/default/carbon.upstart.erb
#!upstart
# http://upstart.ubuntu.com/cookbook/
#
description carbon-<%= @name %>

start on runlevel [2345]
stop on shutdown

expect daemon

respawn
respawn limit 99 5

setuid <%= @user %>

exec <%= @dir %>/bin/carbon-<%= @name %>.py start
