# templates/default/uwsgi.upstart.erb
description "uwsgi graphite instance"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --json <%= node['graphite']['base_dir'] %>/conf/uwsgi-config.json
