# templates/default/uwsgi.upstart.erb
description "uwsgi graphite instance"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --processes <%= node['graphite']['uwsgi']['workers'] %> \
<% if node['graphite']['uwsgi']['carbon'] -%>
--plugins carbon --carbon 127.0.0.1:2003 \
<% end -%>
--pythonpath <%= node['graphite']['base_dir'] %>/lib \
--pythonpath <%= node['graphite']['base_dir'] %>/webapp/graphite \
--wsgi-file <%= node['graphite']['base_dir'] %>/conf/graphite.wsgi.example \
--uid <%= node['graphite']['user_account'] %> --gid <%= node['graphite']['group_account'] %> \
--no-orphans --master \
--buffer-size <%= node['graphite']['buffer_size'] %> \
--procname graphite-web \
--die-on-term \
<% if node['graphite']['uwsgi']['listen_http'] -%>
--http <%= node['graphite']['listen_ip'] unless node['graphite']['listen_ip'].nil? %>:<%= node['graphite']['listen_port'] %>
<% else -%>
--socket <%= node['graphite']['uwsgi']['socket'] %>
<% end -%>
